# host/root/etc/cron.d/hiauntie_backup
PATH=/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

30 3 * * * hiauntie_bot /home/hiauntie_bot/maintenance.hiauntie.com/host/script/backup.sh -d >> /var/log/hiauntie/backup.bot.log 2>&1
