# hiauntie/root/etc/cron.hourly/hiauntie_log_df
#!/bin/bash

set -e

sudo -u hiauntie_bot /home/hiauntie_bot/maintenance.hiauntie.com/hiauntie/script/log_df.sh
